# templates/jupyter_dashboard_server.conf
# jupyter_dashboard_server upstart manager

description "Manage jupyter_dashboard_server with upstart"
start on runlevel [2345]
stop on runlevel [06]

script

  exec jupyter-dashboards-server \
  --KERNEL_GATEWAY_URL='{{ lookup("env", "KERNEL_GATEWAY_URL") }}' \
  --IP=0.0.0.0

  echo 'jupyter_dashboard_server up'

end script

respawn
